# niyam.service
[Unit]
Description=Niyam - CLI Command Governance System
After=network.target

[Service]
Type=simple
User=hangsai
Group=users
WorkingDirectory=/home/hangsai/openclaw-workspace/niyam
ExecStart=/home/hangsai/.nix-profile/bin/node server.js
Environment=PORT=3200
Environment=NODE_ENV=production
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=niyam

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/hangsai/openclaw-workspace/niyam/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
